# docker/app/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:octane]
command=sh -c "exec php /app/artisan octane:start --server=frankenphp --host=0.0.0.0 --port=8000 --workers=${OCTANE_WORKERS:-8} --max-requests=1000 --no-interaction"
autostart=true
autorestart=true
priority=5
user=www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
